# systemd/2dosumi.service
[Unit]
Description=2dosumi second-sleep prevention monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yuki
WorkingDirectory=/home/yuki/2dosumi
EnvironmentFile=-/home/yuki/2dosumi/.env
ExecStart=/home/yuki/2dosumi/.venv/bin/2dosumi run --source hx711
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
